# base/rootfs/overlay/etc/systemd/system/bluetooth.service.d/10-lumelo-rfkill-unblock.conf
[Unit]
After=systemd-rfkill.service systemd-rfkill.socket
Wants=systemd-rfkill.service

[Service]
ExecStartPre=/bin/sh -c 'if command -v rfkill >/dev/null 2>&1; then rfkill unblock bluetooth || true; rfkill unblock all || true; fi'
