# base/rootfs/overlay/etc/systemd/system/bluetooth.service.d/20-lumelo-uart-attach.conf
[Unit]
Requires=lumelo-bluetooth-uart-attach.service
After=lumelo-bluetooth-uart-attach.service
