# base/rootfs/overlay/etc/systemd/system/ssh.service.d/10-lumelo-hostkeys.conf
[Unit]
Requires=lumelo-ssh-hostkeys.service
After=lumelo-ssh-hostkeys.service
