# base/rootfs/overlay/etc/systemd/system/auth-recovery.service
[Unit]
Description=Reset admin authentication from physical recovery media
DefaultDependencies=no
After=local-fs.target
Before=controld.service local-mode.target bridge-mode.target
ConditionPathExists=/usr/libexec/lumelo/auth-recovery

[Service]
Type=oneshot
ExecStart=/usr/libexec/lumelo/auth-recovery

[Install]
WantedBy=local-mode.target bridge-mode.target
