# base/rootfs/overlay/etc/systemd/system/lumelo-media-reconcile.service
[Unit]
Description=Lumelo Removable Media Availability Reconcile
After=local-fs.target
ConditionPathExists=/usr/bin/lumelo-media-import

[Service]
Type=oneshot
ExecStart=/usr/bin/lumelo-media-import reconcile-volumes
